# src/example/stats_update/crontab
* * * * * sleep 10; /usr/bin/php <current_path>/cronjob.php
* * * * * sleep 20; /usr/bin/php <current_path>/cronjob.php
* * * * * sleep 30; /usr/bin/php <current_path>/cronjob.php
* * * * * sleep 40; /usr/bin/php <current_path>/cronjob.php
* * * * * sleep 50; /usr/bin/php <current_path>/cronjob.php